# glowd.service
[Unit]
Description=glowd - MagicHome LED web controller
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/glowd --port 5578
Restart=on-failure
DynamicUser=yes
StateDirectory=glowd
Environment=GLOWD_STATE_DIR=/var/lib/glowd

[Install]
WantedBy=multi-user.target
